# resources/deb/debian/prerm
#!/bin/sh
set -e

case "$1" in
    remove|deconfigure)
        echo "==> Stopping and disabling quadcast2srgb service if running..."
        # Run as the invoking user if possible; fall back silently if no session exists
        systemctl --user disable --now quadcast2srgb.service 2>/dev/null || true
    ;;
    upgrade|failed-upgrade)
    ;;
    *)
        echo "prerm called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
