# resources/systemd/quadcast2srgb.service.in
[Unit]
Description=HyperX QuadCast 2S RGB Controller
Documentation=https://github.com/MuffinMario/quadcast2Srgb
# start after USB subsystem and login session are available
After=local-fs.target systemd-udevd.service

[Service]
Type=notify
# default is to just process the config file that is being shipped
# dirs will be repl during cmake build process
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/qc2srgb --config @CMAKE_INSTALL_FULL_SYSCONFDIR@/quadcast2srgb/config.toml
Restart=on-failure
RestartSec=5s
# watchdog - systemd kills+restarts if no WATCHDOG=1 within 10s
# we set this in the active finder thread to signal we're still alive (maybe instead: inside sender?)
WatchdogSec=10s
# stop request can take a possibly unfortunate amount of time when we have a bunch of devices queued to manage
TimeoutStopSec=30s 
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quadcast2srgb

# no need to expose public tmp dir
PrivateTmp=true 
# everything read only (except hidraw dev explicitly declared in rule)
ProtectSystem=strict 
# only needs read perms for config
ProtectHome=read-only 
# cannot ever elevate
NoNewPrivileges=true 

[Install]
WantedBy=default.target